# deploy/podmanager.service.template
[Unit]
Description=Podcast Manager API
After=network.target postgresql.service

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__/server
ExecStart=__NODE_PATH__ dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=__APP_DIR__/.env

[Install]
WantedBy=multi-user.target
